# services/authorize_key.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

[program:authorize_key]
command=/usr/local/bin/authorize_key
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
# not logging stdout if you want to see logs, remove the following 2 lines
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
